# deploy/va-update-listener.service
[Unit]
Description=Voice Assistant update trigger listener
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'if [ ! -p /tmp/va-update ]; then rm -rf /tmp/va-update && mkfifo /tmp/va-update; fi && chmod 666 /tmp/va-update'
ExecStart=/bin/bash -c 'while true; do cat /tmp/va-update; /opt/voice-assistant/deploy/update.sh; done'
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
